# picore-extension/services/input-handler.service
[Unit]
Description=Input Handler for RP2040 Controller
Documentation=file:///opt/roll-streamer/README.md
After=network.target

[Service]
Type=simple
User=tc
Group=staff
WorkingDirectory=/opt/roll-streamer
ExecStart=/usr/bin/python3 /opt/roll-streamer/scripts/input_handler.py --poll-rate 20
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Resource limits
MemoryLimit=50M
CPUQuota=10%

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
